# files/etc/hotplug.d/net/99-wifi
#!/bin/sh
# JDCloud RE-CS-02 (雅典娜) 专用 WiFi 热插拔修复脚本

[ "$ACTION" = "add" ] || exit 0

WIRELESS_CFG="/etc/config/wireless"
PCIE_PATH="soc@0/20000000.pcie/pci0000:00/0000:00:00.0/0000:01:00.0"

# 仅当配置文件不存在，或存在但缺少正确的 radio2 PCIe 路径时触发修复
if [ ! -f "$WIRELESS_CFG" ] || ! grep -q "$PCIE_PATH" "$WIRELESS_CFG"; then
    logger -t wifi-fix "RE-CS-02: Wireless config missing or radio2 path incorrect, regenerating..."
    
    # 备份旧配置（如果存在）
    [ -f "$WIRELESS_CFG" ] && cp "$WIRELESS_CFG" "${WIRELESS_CFG}.bak.$(date +%s)"
    
    # 调用系统自带生成器重新生成基础配置
    wifi config > /dev/null 2>&1
    
    # 强制修正 radio2 的路径（防止自动生成器写错）
    if grep -q "radio2" "$WIRELESS_CFG"; then
        uci set wireless.radio2.path="$PCIE_PATH"
        uci set wireless.radio2.band='5g'
        uci set wireless.radio2.channel='36'
        uci set wireless.radio2.htmode='HE80'
        uci commit wireless
        logger -t wifi-fix "RE-CS-02: radio2 PCIe path fixed successfully."
    else
        logger -t wifi-fix "RE-CS-02: WARNING! radio2 not detected by wifi config. Check ath11k_pci driver and ART partition."
    fi
fi
