# files/etc/init.d/dnsmasq-delay
#!/bin/sh /etc/rc.common
START=19
USE_PROCD=1

start_service() {
    # 等待 br-lan 接口就绪（最多15秒）
    local i=0
    while [ $i -lt 15 ]; do
        if [ -e /sys/class/net/br-lan/operstate ] && \
           [ "$(cat /sys/class/net/br-lan/operstate 2>/dev/null)" = "up" ]; then
            break
        fi
        sleep 1
        i=$((i+1))
    done

    # 直接重载 dnsmasq 配置（更可靠）
    /etc/init.d/dnsmasq reload
    logger -t DIY "dnsmasq reloaded after $i seconds"
}

stop_service() {
    :
}
